# scripts/launchd/com.synapse.weekly-release.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Synapse weekly-release — consolidates the week into ONE PR on Fridays 18:03
     (driver §4.5 / §6: cron "3 18 * * 5"). Never merges (owner gate).
     Placeholders __HOME__ / __REPO__ are substituted by scripts/launchd/install.sh. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.synapse.weekly-release</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO__/scripts/loop-run.sh</string>
    <string>release</string>
  </array>

  <!-- Friday 18:03 (Weekday 5 = Friday). Closes the Sat→Fri cycle. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>5</integer>
    <key>Hour</key><integer>18</integer>
    <key>Minute</key><integer>3</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>__HOME__</string>
    <key>PATH</key><string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- Owner authorized full no-input autonomy 2026-07-12. Revert to acceptEdits to disarm. -->
    <key>LOOP_PERMISSION_MODE</key><string>bypassPermissions</string>
  </dict>

  <key>StandardOutPath</key><string>__HOME__/Library/Logs/synapse-loop/launchd-release.out.log</string>
  <key>StandardErrorPath</key><string>__HOME__/Library/Logs/synapse-loop/launchd-release.err.log</string>

  <key>ProcessType</key><string>Background</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
